# dist/static/auto_deploy_scripts/auto_deploy.service
[Unit]
Description=VM Auto Deploy Server

[Service]
Type=simple
WorkingDirectory=/home/{{user}}
Environment="PYTHONPATH=/home/{{user}}"
User={{user}}
ExecStart=/home/{{user}}/venv/bin/python3 /home/{{user}}/venv/bin/gunicorn auto_deploy:app
TasksMax=50
MemoryMax=52428800

[Install]
WantedBy=multi-user.target
